# vm/openstack/clustercheck/systemd/mysqlchk@.service
[Unit]
Description=Percona XtraDB Cluster node check service

[Service]
ExecStart=-/usr/local/bin/clustercheck
StandardInput=socket